# docker/supervisord.conf
[supervisord]
nodaemon=true
[program:node]
command=adonis kue:listen
autostart=true
autorestart=true